# Linux-Hedge_Project/src/MongodbServer/etc/systemd/system/mongod.service.d/limits.conf
[Service]
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitMEMLOCK=infinity
LimitNOFILE=65535
LimitNPROC=65535
TasksMax=infinity